# conf/startup.conf
[program:btsync]
command=btsync /btsync/btsync.conf --nodaemon
redirect_stderr=true
stdout_logfile=/tmp/btsync.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3